# misc/github-weather.crontab
WORKDIR=/home/pi/.local/lib/github-weather

*/10 * * * * $WORKDIR/github-weather.linux-arm7 -configuration $WORKDIR/config.yaml 2>> /tmp/github-weather.log
