# misc/github-weather.crontab
*/10 * * * * /home/pi/github-weather.linux-arm7
